# configs/jabali-terminal.service
[Unit]
Description=Jabali Terminal — WebSocket-to-PTY daemon for browser root shell
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/jabali-terminal/venv/bin/python -m jabali_terminal.server
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jabali-terminal

# Security hardening
User=root
Group=root
# NOTE (SEC-REV-9): NoNewPrivileges=no required so daemon can exec /bin/bash under PTY.
# PTY bridge must be root to assign pty-master to www-data group (0660 perms).
NoNewPrivileges=no
# NOTE (SEC-REV-9): ProtectSystem=false required so daemon can read/write session logs
# and nonce store in /var/lib/jabali-terminal/ and /var/log/jabali-terminal/.
# Daemon does NOT write to /etc or /usr; logs and nonces are owned root:root 0600.
ProtectSystem=false
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=no

# Runtime directories (auto-created if absent)
RuntimeDirectory=jabali-terminal
RuntimeDirectoryMode=0755
StateDirectory=jabali-terminal
StateDirectoryMode=0755
LogsDirectory=jabali-terminal
LogsDirectoryMode=0700

[Install]
WantedBy=multi-user.target
